# etc/profile.d/emacs.sh
ALTERNATE_EDITOR=/usr/local/bin/emacs
# ALTERNATE_EDITOR=''
export ALTERNATE_EDITOR
export EDITOR=/usr/local/bin/emacsclient
